# packaging/vocala.service
[Unit]
Description=Vocala Voice Chat Server
After=network.target

[Service]
Type=simple
User=vocala
Group=vocala
ExecStart=/usr/bin/vocala -config /etc/vocala/config.yaml
WorkingDirectory=/usr/share/vocala
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/vocala

[Install]
WantedBy=multi-user.target
